# modules.d/95zfcp/56-zfcp.rules
KERNEL=="zfcp_cfdc", RUN+="/sbin/zfcpconf.sh"
