# test/kdump_tests_lib/kexec-kdump-test/kexec-kdump-test.service
[Unit]
Description=Kexec Kdump Test Service

[Service]
ExecStart=/kexec-kdump-test/init.sh
Type=idle

[Install]
WantedBy=multi-user.target
