# modules.d/50plymouth/52-plymouth-start.rules
ACTION=="add", SUBSYSTEM=="graphics", KERNEL=="fb[0-9]*", RUN+="/bin/plymouth --show-splash"
