# modules.d/40nfsroot/60-nfsroot.rules
ACTION=="online", SUBSYSTEM=="net", RUN+="/sbin/nfsroot $env{INTERFACE}"
