# modules.d/35network-manager/dracut-wait-nm.path
[Unit]
Description=Wait NetworkManager startup complete

[Path]
PathExists=/run/NetworkManager/startup-complete

[Install]
WantedBy=sysinit.target
