# modules.d/10redhat-i18n/10-console.rules
# Console initialization - keyboard, font, etc.
KERNEL=="tty0",		RUN+="/lib/udev/console_init %k"
