# modules.d/01systemd-initrd/initrd-keep.service
[Unit]
Description=Hang in initrd.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=
